# systemd/fcc-refresh.timer
[Unit]
Description=Weekly FCC Database Refresh (Sunday 03:00)
Requires=fcc-refresh.service

[Timer]
OnCalendar=Sun *-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
