# systemd/repeater-refresh.service
[Unit]
Description=FieldComms Repeater Database Refresh (RepeaterBook)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=fieldcomms
Group=fieldcomms
WorkingDirectory=/opt/fieldcomms
# REPEATERBOOK_TOKEN is read from the environment file below.
# Create /opt/fieldcomms/data/repeaterbook.env with a line:
#   REPEATERBOOK_TOKEN=your-issued-token
EnvironmentFile=-/opt/fieldcomms/data/repeaterbook.env
ExecStart=/opt/fieldcomms/venv/bin/python /opt/fieldcomms/python/fetch_repeaters.py --states IL,WI,IN,IA --bands 2m,70cm
StandardOutput=journal
StandardError=journal
SyslogIdentifier=repeater-refresh
